# goodbyedpi.service
[Unit]
Description=GoodbyeDPI-Linux DPI circumvention service
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/bash /opt/goodbyedpi/scripts/setup_iptables.sh
ExecStart=/opt/goodbyedpi/src/goodbyedpi
ExecStopPost=/bin/bash /opt/goodbyedpi/scripts/remove_iptables.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
